# build/service/sshd/run
#!/bin/bash
set -e

if [ ! -e /var/log/auth ]; then
  touch /var/log/auth
fi

tail -F /var/log/auth &

exec /usr/sbin/sshd -D